# salt/base/k3s/k3s.override.service
[Service]
ExecStop=/usr/bin/k3s-killall
